# silver-accounting-bot/deploy/systemd/worker.service
[Unit]
Description=Silver Accounting Bot ARQ Worker
After=redis.service postgresql.service bot.service
Requires=redis.service postgresql.service

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/opt/silver-accounting-bot
EnvironmentFile=/opt/silver-accounting-bot/.env
ExecStart=/opt/silver-accounting-bot/.venv/bin/python -m infrastructure.tasks.arq_worker
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=silver-worker

[Install]
WantedBy=multi-user.target
